# etc/NetworkManager/dispatcher.d/99-dwmblocks-update
#!/usr/bin/env bash

pkill -USR1 -f "dwmblocks"

